# vsplatego.conf
[program:vsplatego]
command=/bin/sh /var/www/html/vsplatego.sh

